# setup/supervisor/rmotr-execute.ini
[program:rmotr-execute]
command=authbind %(ENV_HOME)s/.virtualenvs/rmotr-execute/bin/gunicorn -w 4 -b 0.0.0.0:4000 api.app:app
environment=PYTHONPATH="%(ENV_HOME)s/rmotr-execute"
numprocs=1
autostart=true
stderr_logfile=/var/log/supervisor/rmotr-execute.stderr.log
stdout_logfile=/var/log/supervisor/rmotr-execute.stdout.log
autorestart=unexpected
process_name=rmotr-execute
